# C04_Realtime-Recommendation/src/module2/task2/stage2/flink.cron
* * * * * python /root/Downloads/query_generator.py /tmp/flink-datasets/$(date + "\%Y-\%m-\%d-\%H-\%M-\%S").log 1000
